# templates/etc/supervisor/conf.d/celery.conf.j2
; << Ansible managed file >>

[program:{{item.key}}]
{% for k,v in item.value.iteritems() %}
{{k}} = {{v}}
{% endfor %}

